# tempcheck.timer
[Unit]
Description=Run temperature check periodically

[Timer]
# Primo avvio 30s dopo il boot, poi ogni 30s.
OnBootSec=30
OnUnitActiveSec=30
# Niente AccuracySec=1ms: lasciare a systemd un minimo di margine evita
# wakeup continui della CPU. 5s di tolleranza vanno benissimo per un alert.
AccuracySec=5s
# Recupera un'esecuzione persa se il sistema era spento/sospeso.
Persistent=true

[Install]
WantedBy=timers.target
