# tempcheck.service
[Unit]
Description=Check temperature and notify via Telegram
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot

# --- User that runs the script ---
# Default: root (reads sensors and writes the state file without issues).
# To use an unprivileged user:
#   1) change User= and Group= below
#   2) make sure it can read 'sensors' (coretemp is usually world-readable)
#   3) make sure ALERT_FILE lives in a directory writable by that user
User=root
#Group=root

# Path of the script (must match where you copy it).
ExecStart=/root/alertemp.sh

# Credentials are read from /etc/alertemp.env by the script itself.
# Alternatively you can pass them here (not recommended, they end up in logs):
#Environment=THRESHOLD=80
#EnvironmentFile=-/etc/alertemp.env

# Minimal hardening (optional, uncomment to restrict)
#NoNewPrivileges=true
#ProtectSystem=strict
#ProtectHome=true
#ReadWritePaths=/tmp
