# tempcheck.timer
[Unit]
Description=Run temperature check periodically

[Timer]
# First run 30s after boot, then every 30s.
OnBootSec=30
OnUnitActiveSec=30
# No AccuracySec=1ms: giving systemd some margin avoids constant CPU
# wakeups. 5s of tolerance is perfectly fine for an alert.
AccuracySec=5s
# Recover a missed run if the system was powered off/suspended.
Persistent=true

[Install]
WantedBy=timers.target
