# old/tempcheck.service
[Unit]
Description=Check temperature and notify
[Service]
ExecStart=/root/alertemp.sh
